# 81-usbserial.rules
ACTION=="add", SUBSYSTEM=="tty", KERNEL=="ttyUSB0", RUN+="/bin/systemctl start serialdump"
